# templates/etc/systemd/system/docker.service.d/docker.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon {{ zenoss_docker_opts }}
